# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=bash -c "cd /app && python scripts/init_db.py && alembic upgrade head && uvicorn src.main:app --host 0.0.0.0 --port 8000"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="sqlite+aiosqlite:////app/data/mcparr.db",REDIS_URL="redis://localhost:6379",LOG_LEVEL="INFO"

# MCP Server - Disabled (MCP uses stdio for Claude Desktop, not HTTP)
# For Open WebUI integration, use the OpenAPI tools endpoint: /tools/openapi.json
# [program:mcp-server]
# command=python -m src.mcp.main
# directory=/app
# autostart=false
# autorestart=false
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
